# systemd files/biometric-server.service
[Unit]
Description=Biometric WebSocket Server
After=network-online.target       
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/biometric-server
ExecStart=/usr/bin/python3 /opt/biometric-server/biometric_to_server.py
Restart=always
RestartSec=5                     
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/biometric-server
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
